# SMPPServer/files/systemd/SMPPAuthenticator@.service
# /etc/systemd/system/SMPPAuthenticator@.service
#
# Template unit for the SMPP Authenticator.
# Started by simple_smpp_server after a successful TCP handshake.
#
# The instance parameter (%i) carries the session ID, e.g.:
#   systemctl start SMPPAuthenticator@0195f3a7b2c1-a3f8d901c4b72e5d.service
#
# Install:
#   sudo cp SMPPAuthenticator@.service /etc/systemd/system/
#   sudo systemctl daemon-reload

[Unit]
Description=SMPP Authenticator – session %i
Documentation=man:smpp_authenticator(8)
After=network.target

[Service]
Type=simple

# Replace /usr/local/bin/smpp_authenticator with the actual binary path.
ExecStart=/usr/local/bin/smpp_authenticator --session-id %i

# Pass the session ID via environment as well, for libraries that read it.
Environment="SMPP_SESSION_ID=%i"

Restart=on-failure
RestartSec=3s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/smpp_authenticator

# Logging – journald captures stdout/stderr automatically.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smpp_auth_%i

[Install]
# Typically not enabled statically; started on-demand by simple_smpp_server.
WantedBy=multi-user.target
