# Docker/smsc-dev/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# D-Bus system bus daemon (for development)
[program:dbus]
command=/usr/bin/dbus-daemon --config-file=/etc/dbus-1/system.conf --nofork
autostart=true
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus_err.log
startsecs=2
stopasgroup=true

# SSH daemon
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
priority=200
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd_err.log
startsecs=2
stopasgroup=true

# Unix socket for supervisor
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

# RPC interface
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# CLI for supervisorctl
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
